# rpi-climate-monitor-measurer/crontab
* * * * * root /data/rpi-climate-monitor/performMeasurement
