# crontab
* * * * * /data/rpi-climate-monitor/performMeasurement
